# infrastructure/environment-importer/systemd/orchidapp-environment-importer.service
[Unit]
Description=OrchidApp environment email importer
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=raymond-23
Group=raymond-23
EnvironmentFile=/etc/orchidapp/environment-importer.env
ExecStart=/opt/orchidapp/infrastructure/environment-importer/.venv/bin/python /opt/orchidapp/infrastructure/environment-importer/import-environment-email.py --config /etc/orchidapp/environment-importer.json